# bitcoind.service.tmpl
# vim: ft=systemd
# .service files are generated from .service.tmpl files
[Unit]
Description=bitcoind-__NETWORK__

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=20m
ExecStartPre=/bin/sh -c "IMAGE=olalond3/coreos-bitcoind && docker pull $IMAGE"
ExecStartPre=/bin/sh -c "docker inspect bitcoind-__NETWORK__ >/dev/null 2>&1 \
  && docker rm -f bitcoind-__NETWORK__ || true"
ExecStart=/bin/sh -c "IMAGE=olalond3/coreos-bitcoind \
  && docker run --name bitcoind-__NETWORK__ --rm \
    --net=host \
    -e HOST=${COREOS_PRIVATE_IPV4} \
    -e NETWORK=__NETWORK__ \
    $IMAGE"
ExecStop=-/usr/bin/docker stop bitcoind-__NETWORK__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

[X-Fleet]
Conflicts=bitcoind-__NETWORK__.service
